# deploy/conf/supervisord-mysql.conf
[program:mysqld]
command=/start-mysql.sh
numprocs=1
autostart=true
autorestart=true
