# deploy/conf/supervisord-apache2.conf
[program:apache2]
command=/start-apache2.sh
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
