# deploy/conf/supervisord-apache2.conf
[program:apache2]
command=/start-apache2.sh
numprocs=1
autostart=true
autorestart=true
