# deploy/conf/supervisord-logging.conf
[program:logging]
command=/start-logging.sh
numprocs=1
autostart=true
autorestart=true
#stdout_logfile=/proc/1/fd/1
#stderr_logfile=/proc/1/fd/1

